# scripts/service/cerbere.service
[Unit]
Description=Cerbere RADIUS service
